# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:xvfb]
command=/usr/bin/Xvfb :%(ENV_DISPLAY_NUM)s -screen 0 1920x1080x24 -ac
autorestart=true
priority=10

[program:x11vnc]
command=/usr/bin/x11vnc -display :%(ENV_DISPLAY_NUM)s -rfbport 5999 -nopw -forever -shared -xkb
autorestart=true
priority=20
startsecs=2

[program:websockify]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5999
autorestart=true
priority=30
startsecs=2

[program:crawler-api]
command=python -m uvicorn api_server:app --host 0.0.0.0 --port %(ENV_CRAWLER_API_PORT)s
directory=/opt/crawler
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s",CRAWLER_WORKING_DIR="/opt/crawler",CRAWLER_DATA_DIR="/opt/crawler/output",CRAWLER_VENV_PYTHON="/usr/local/bin/python",MAX_DISPLAY_SESSIONS="%(ENV_MAX_DISPLAY_SESSIONS)s"
autorestart=true
priority=40
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:crawler-data]
command=python data_server.py
directory=/opt/crawler
environment=CRAWLER_DATA_DIR="/opt/crawler/output",CRAWLER_DATA_PORT="%(ENV_CRAWLER_DATA_PORT)s"
autorestart=true
priority=40
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
